# config/lswitch.service
[Unit]
Description=LSwitch - Layout Switcher (evdev)
Documentation=https://github.com/yourusername/lswitch
After=network.target multi-user.target

[Service]
Type=simple
# Prefer running the package as a module and ensure the installed
# package dir is available via PYTHONPATH so imports resolve correctly.
ExecStart=/usr/local/bin/lswitch-run
# NOTE: `lswitch-run` will ensure the installed package dir is available on PYTHONPATH and
# avoid name-shadowing issues that can occur when a script shares the package name.
Restart=on-failure
RestartSec=5s
WorkingDirectory=/etc/lswitch
Environment="PYTHONUNBUFFERED=1"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/anton/.Xauthority"

# Безопасность
# NoNewPrivileges=true
# PrivateTmp=true

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lswitch

[Install]
WantedBy=multi-user.target
