# archive/config/lswitch.service
[Unit]
Description=LSwitch - Layout Switcher (evdev)
Documentation=https://github.com/kabobik/lswitch
After=graphical-session.target

[Service]
Type=simple
# Используем sh -c для поиска lswitch в PATH
ExecStart=/bin/sh -c 'exec $(command -v lswitch || echo /usr/local/bin/lswitch)'
Restart=on-failure
RestartSec=5s
Environment="PYTHONUNBUFFERED=1"
# DISPLAY наследуется из user session
Environment="XAUTHORITY=%h/.Xauthority"

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lswitch

[Install]
WantedBy=default.target
